# server/etc/network/if-up.d/00-upstart
#!/bin/sh -e

MARK_DEV_PREFIX="/run/network/ifup."

[ ! -e "$MARK_DEV_PREFIX$IFACE" ] || exit 0

. /etc/network/if-up.d/upstart

